# scripts/rpi/systemd/safecare-docker.service
[Unit]
Description=SafeCare Docker Compose Services
ConditionPathExists=/opt/safecare/.provisioned
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/safecare/docker
# docker-up.sh reads .variant to select the right compose profile
ExecStart=/opt/safecare/scripts/rpi/docker-up.sh up -d
ExecStop=/opt/safecare/scripts/rpi/docker-up.sh down

[Install]
WantedBy=multi-user.target
